# config/includes.chroot/etc/systemd/system/fellowshipos-setup.service
[Unit]
Description=FellowshipOS First Boot Setup Wizard
After=network.target
Before=getty.target
Before=getty@tty1.service
ConditionPathExists=!/etc/fellowshipos/.setup-complete

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fellowshipos-setup
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
